# config/music/.config/autostart/mpdas.desktop
[Desktop Entry]
Type=Application
Exec=mpdas
X-GNOME-Autostart-enabled=true
Name=mpdas
Comment=Scrobbling support for the MPD
