# config/music/.config/autostart/mpdris.desktop
[Desktop Entry]
Type=Application
Exec=mpd-mpris
X-GNOME-Autostart-enabled=true
Name=MPDRIS
Comment=MPRIS support for the MPD
